# one-thing-done.service
[Unit]
Description=One Thing Done - HTTP Command Executor
Documentation=https://github.com/rmrf/one-thing-done
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/one-thing-done --config /usr/local/etc/one-thing-done-config.toml
Restart=on-failure
RestartSec=5

# 允许连接到 systemd D-Bus
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/systemd/private"

# 基本安全设置（不影响 systemctl）
NoNewPrivileges=false
ProtectSystem=no
ProtectHome=no

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=one-thing-done

[Install]
WantedBy=multi-user.target
